# root/etc/udev/old/ps3controller.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="054c", ATTR{idProduct}=="0268", RUN+="/usr/local/bin/udevrun.sh /usr/local/bin/sixpair.sh"

