# root/etc/udev/rules.d/91-btremove.rules
ACTION=="remove", KERNEL=="hci0", RUN+="/usr/local/bin/udevrun.sh btstop.sh"
