# root/etc/udev/rules.d/92-btnetadd.rules
ACTION=="add", KERNEL=="bnep0", RUN+="/usr/local/bin/udevrun.sh /usr/local/bin/btnetstart.sh"
